# services/cron/rootfs/var/spool/cron/crontabs/app
# Put the cron jobs here.
0 0 * * * cd /app; sh backup_treemap_database.sh
