# services/nginx/rootfs/var/spool/cron/crontabs/root
# thing here yet
0 0 * * 1 cd /app/ssl; ./update-all
